# dht22.service
[Unit]
Description=Temperature and Humidity Sensor

[Service]
Type=idle
ExecStart=/usr/bin/python3 /home/toby/grayes-monitors/dht22.py
Restart=always
User=toby

[Install]
WantedBy=multi-user.target
